# root/etc/s6-overlay/s6-rc.d/svc-wgdashboard/run
#!/usr/bin/with-contenv bash

cd /wgdashboard/src
./wgd.sh debug
